# systemd/spotifyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/spotifyd --no-daemon --config-path /home/aignas/.config/spotifyd/spotify.conf
