# spotifyd/spotifyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/spotifyd --no-daemon --config-path="${HOME}/.config/spotifyd/spotify.conf" --cache-path="${HOME}/.cache/spotify" --use-mpris=true --backend=pulseaudio
